# dns-health-check/dns-health-check.timer
[Unit]
Description=Periodic DNS health check
Documentation=man:systemd-resolved.service(8)

[Timer]
# Run shortly after login, then every 5 minutes.
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=dns-health-check.service
Persistent=true

[Install]
WantedBy=timers.target
