# dns-health-check/dns-health-check.service
[Unit]
Description=DNS health check (systemd-resolved degradation monitor)
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dns-health-check
Nice=10
# Allow notify-send to reach the user's bus.
PassEnvironment=DBUS_SESSION_BUS_ADDRESS WAYLAND_DISPLAY DISPLAY XDG_RUNTIME_DIR

[Install]
WantedBy=default.target
