# aerc-mail/.config/systemd/user/mbsync.service
[Unit]
Description=Sync all mail accounts into the local Maildir with mbsync (isync)
Documentation=man:mbsync(1)

[Service]
Type=oneshot
# Sync every channel defined in ~/.config/isyncrc. Auth uses OAuth2: mbsync runs
# `oama access <addr>` per account, which prints a fresh token (oama refreshes in
# the background and keeps tokens GPG-encrypted). aerc watches the Maildir and
# refreshes on its own when this runs.
ExecStart=/usr/bin/mbsync -a
# A failed run (offline, expired/unauthorized token) shouldn't be noisy;
# the timer just tries again next interval.
SuccessExitStatus=0 1
